# mclient
#!/sbin/runscript
# Copyright 2014 Seth Walsh
# Distributed under the terms of the GNU General Public License v2

description="Monitoring-Client user and program statistical gathering service"

# Dependencies of Monitoring-Client
depend() {
  need net
  need localmount
}

start() {
  ebegin "Starting Monitoring-Client"
  start-stop-daemon  --background --start --exec /opt/monitoring/monitoring-client &>/dev/null --pidfile /var/run/monitoring-client.pid --make-pidfile --name monitoring-client
  eend $?
}

stop() {
  ebegin "Stopping Monitoring-Client"
  start-stop-daemon --stop --exec /opt/monitoring/monitoring-client --pidfile /var/run/monitoring-client.pid
  eend $?
}
